# src/etc/systemd/system/dhcpcd5.dhcpcd@.service.d/override.conf
[Service]
ExecStart=/usr/sbin/dhcpcd -6 -c /usr/libexec/vyos/system/vyos-dhcpcd-script.py -f /run/dhcpcd/dhcpcd.%I.conf -q -w %I
